# env_tg_bot.service
[Unit]
Description=env tg bot

[Service]
TimeoutStartSec=0
RestartSec=10
Restart=always

KillSignal=SIGTERM
PIDFile=/tmp/env_tg_bot.pid
Type=simple

User=USER

# Main process
ExecStart=/bin/bash PATH_TO_RUN_SH
ExecStop=/bin/kill -TERM $(cat /tmp/env_tg_bot.pid)

[Install]
WantedBy=multi-user.target